# database/mysql/db-backup/sdb03_mysql_dump.timer
[Unit]
Description=每天凌晨2点执行

[Timer]
OnCalendar=*-*-* 2:0:0
Persistent=true

[Install]
WantedBy=multi-user.target
