# database/mysql/db-backup/sdb03_mysql_dump.service
[Unit]
Description=MySQL数据库备份

[Service]
Type=simple
WorkingDirectory=/data/mnt
ExecStart=/data/mnt/backup.sh
User=mysql
Group=mysql